# 11-detection-engineering-sigma/04-scheduled-task-persistence/sigma-rule.yml
title: Suspicious Scheduled Task Creation via PowerShell
id: 7d2c5b41-5f67-4f2c-8e2a-9d4a8f2e1c11
status: experimental
description: Detects PowerShell spawning schtasks.exe to create a scheduled task.

logsource:
  product: windows
  category: process_creation

detection:
  selection_parent:
    ParentImage|endswith: '\powershell.exe'
  selection_child:
    Image|endswith: '\schtasks.exe'
  selection_command:
    CommandLine|contains: '/Create'
  condition: selection_parent and selection_child and selection_command

falsepositives:
  - Administrative task creation
  - Software deployment activities

level: high

tags:
  - attack.persistence
  - attack.t1053.005
  - attack.execution
  - attack.t1059.001
  - attack.defense_evasion
  - attack.t1027